# target/linux/ar71xx/base-files/etc/uci-defaults/01-factory
#!/bin/sh

change_default_network() {
	local ret=$(uci -q get network.MOBILE)	

	if [ "$ret" == "interface" ]; then
		uci delete network.MOBILE
		uci set network.lan.ifname=eth0
		uci set network.lan.ipaddr=192.168.10.1
		uci set network.wan.ifname=eth1
		uci commit network

		uci delete mwan3.MOBILE
		uci delete mwan3.MOBILE_m
		uci delete mwan3.4G_only
		uci delete mwan3.ALL
		uci commit mwan3

		uci set system.@system[0].hostname=OutdoorExtender
		uci commit system
	fi
}


distname=$(grep DISTRIB_ID /etc/openwrt_release | awk -F= '{print $2}' | sed s/\'//g)

serial_number=$(fw_printenv serial_number|awk -F= '{print $2}')

[ "y$serial_number" == "y" ] && {
	local mac=$(ip link show eth0 | grep ether | awk '{ print $2 }' | sed s/\://g)
	local random=$(head -n 10 /dev/urandom | md5sum )
	random=$(echo ${random:0:4})
	serial_number=ES${random}${mac:8:12}
	serial_number=$(echo $serial_number | tr 'a-z' 'A-Z')
	fw_setenv serial_number $serial_number
}

uci set productinfo.hardware.serial_number=$serial_number
uci set productinfo.hardware.product_name=$distname
uci commit productinfo


[ "$distname" == "EZR13" ] && change_default_network

exit 0
