# target/linux/ramips/base-files/etc/uci-defaults/20-wireless
#!/bin/sh

[ -f "/etc/config/wireless" ] && return

wifi config > /etc/config/wireless

exit 0
