# luci/applications/luci-app-wifi-mtk/root/etc/init.d/wireless
#!/bin/sh /etc/rc.common
# xiaguohua
# 2014/6/11

START=70

boot() {
	wifi up
}

start() {
	return
}

stop() {
	return
}

restart() {
	wifi up
}
