# luci/applications/luci-app-wifi-mtk/root/etc/uci-defaults/luci-wifi
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@wireless[-1]
	add ucitrack wireless
	set ucitrack.@wireless[-1].init=wireless
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
